# daemon-example/daemon-example.service
[Unit]
Description=Daemon Example

[Service]
Type=forking
ExecStart=/usr/local/bin/daemon-example

[Install]
WantedBy=multi-user.target
